# systemd/webappenc-initializer.service.in
[Unit]
Description=Import APP DEKs for Preloaded Web App Decryption
DefaultDependencies=no
Requires=central-key-manager.service
After=central-key-manager.service

[Service]
User=security_fw
Group=security_fw
Type=oneshot
ExecStart=/usr/bin/wae_initializer
SmackProcessLabel=System

[Install]
WantedBy=multi-user.target
